# notif-server.service
[Unit]
Description=Notification Server
After=network.target

[Service]
Type=simple
# Change this to match your Raspberry Pi username (default is 'pi')
User=pi
# %h expands to the home directory of the user above
WorkingDirectory=%h/notif-server
ExecStart=%h/notif-server/target/release/notif-server
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=notif-server

[Install]
WantedBy=multi-user.target
